# openclaw-bridge.service
[Unit]
Description=OpenClaw ESP32 Voice Relay Bridge
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/openclaw-voice-server
ExecStart=/usr/bin/node bridge.js
Restart=on-failure
RestartSec=5
# Adjust environment variables accordingly, or use an EnvironmentFile
Environment=NODE_ENV=production
Environment=PORT=18790
# Environment=ELEVENLABS_API_KEY=YOUR_API_KEY
# Environment=ELEVENLABS_VOICE_ID=YOUR_VOICE_ID

[Install]
WantedBy=multi-user.target
